# deploy/bepository.container
[Unit]
Description=bepository BEP bridge
After=network-online.target
Wants=network-online.target
Conflicts=sleep.target
Before=sleep.target

[Container]
Image=ghcr.io/unbrice/bepository:latest
AutoUpdate=registry

Volume=%S/bepository:/data:z
Volume=%C/bepository:/var/cache/bepository:z

# BEPOSITORY_PORT is expanded by systemd at runtime for both port publishing
# and the application's listen address.
PublishPort=${BEPOSITORY_PORT:-22001}:${BEPOSITORY_PORT:-22001}
ReadOnly=true
Tmpfs=/tmp

# Passes matched variables from the systemd service environment into the
# container. This includes cloud credentials and all BEPOSITORY_* tuning.
Environment=AWS_* GOOGLE_* BEPOSITORY_* CACHE_DIRECTORY

Exec=serve --listen 0.0.0.0:${BEPOSITORY_PORT:-22001}

[Service]
Restart=on-failure
RestartSec=10s
CacheDirectory=bepository
StateDirectory=bepository
EnvironmentFile=/etc/bepository/env
# For credentials on deployments that need eg: per-machine creds.
EnvironmentFile=-/etc/bepository/credentials

[Install]
WantedBy=default.target
